# overlays/finder/etc/systemd/system/finder.service
[Unit]
Description=UDP Finder
Wants=network.target

[Service]
ExecStart=/opt/finder/daemon.py
Restart=always

[Install]
WantedBy=basic.target
